# haproxy/rootfs/etc/services.d/haproxy/run
#!/usr/bin/env sh
exec haproxy -p /var/run/haproxy.pid -f /etc/haproxy/haproxy.cfg